# runit/log/run
#!/bin/sh

mkdir -p /mnt/log/gumshoedb/
exec svlogd -ttt /mnt/log/gumshoedb/
